# config/includes.chroot/etc/skel/.config/autostart/chanx-welcome.desktop
[Desktop Entry]
Type=Application
Name=chanxOS Welcome Center
Exec=chanx-welcome
X-GNOME-Autostart-enabled=true
