# config/includes.chroot/etc/skel/.config/autostart/picom.desktop
[Desktop Entry]
Type=Application
Name=Picom Compositor
Comment=chanxOS liquid glass compositor
Exec=bash -lc 'sleep 2; chanx-compositor'
X-GNOME-Autostart-enabled=true
OnlyShowIn=XFCE;
Terminal=false
Hidden=false
