# config/includes.chroot/etc/skel/.config/autostart/plank.desktop
[Desktop Entry]
Type=Application
Name=Plank
Exec=plank
X-GNOME-Autostart-enabled=true
NoDisplay=false
Terminal=false
Hidden=false
